# slswitch/slswitch.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/opt/lxdeapps/slswitch/slswitch"
